# web-interface/tmux-web.service
[Unit]
Description=Tmux WebSocket Web Interface
After=network.target

[Service]
Type=simple
User=jclee
WorkingDirectory=/home/jclee/app/tmux/web-interface
ExecStart=/usr/bin/node /home/jclee/app/tmux/web-interface/server.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=PORT=3030

# Logging
StandardOutput=append:/var/log/tmux-web.log
StandardError=append:/var/log/tmux-web-error.log

[Install]
WantedBy=multi-user.target
